# etc/cron.daily/00-aws-ec2-smart-snapshot.cron
#!/bin/sh

/usr/local/sbin/aws-ec2-smart-snapshot >> /var/log/aws-ec2-smart-snapshot.log 2>&1
